# extras/cron.d/ackstorm-sync
# Start on boot
@reboot	root	/usr/local/ackstorm/sync/bin/ackstorm-sync start
