# extras/cron.d/ackstorm-sync-master
# Restart every 4 hours to flush inotify handlers
00 */4	* * *	root	/usr/local/ackstorm/sync/bin/ackstorm-sync restart


# Clean data files (10080: 1 week)
00 *	* * *	root	find /usr/local/ackstorm/sync/data -depth -mindepth 1 -maxdepth 1 -ignore_readdir_race -type f -cmin 10080 -delete

# In sync file touch
*  *    * * *   root	touch /tmp/sync-client.done
